# SomfyCameraStream/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
priority=1

[program:mediamtx]
command=/usr/bin/mediamtx /etc/mediamtx/mediamtx.yml
autostart=true
autorestart=true
